# deploy/radar.service
[Unit]
Description=NSE Trade Opportunity Radar (multi-stock scanner, Telegram alerts, edge tracking)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=emcure
Group=emcure
WorkingDirectory=/opt/emcure
ExecStart=/opt/emcure/venv/bin/python3 -m apps.radar_headless
Restart=always
RestartSec=15
EnvironmentFile=/opt/emcure/.env
Environment=PYTHONUNBUFFERED=1
Environment=LOKY_MAX_CPU_COUNT=1
Environment=OMP_NUM_THREADS=1
# Optional: override the default 5-minute scan interval
# Environment=RADAR_REFRESH_SECONDS=300
StandardOutput=append:/var/log/emcure/radar.log
StandardError=append:/var/log/emcure/radar.err

[Install]
WantedBy=multi-user.target
